# .config/systemd/user/https.socket
[Socket]
ListenStream=0.0.0.0:443
FileDescriptorName=websecure
Service=traefik.service

[Install]
WantedBy=sockets.target
